# files/koha-hourly.cron
#!/usr/bin/env bash

/usr/share/koha/bin/cronjobs/holds/build_holds_queue.pl
